# supervisord.conf
[supervisord]
nodaemon=true

[program:fauxton]
command=sh /docker-entrypoint.sh
autostart=true
autorestart=true